# vps-agent/writgo-agent.service
[Unit]
Description=Writgo VPS Agent - Automated WordPress Publishing
After=network.target redis.service

[Service]
Type=simple
User=root
WorkingDirectory=/opt/writgo-vps/vps-agent
ExecStart=/usr/bin/node dist/index.js
Restart=always
RestartSec=10
StandardOutput=journal
StandardError=journal

# Environment
Environment="NODE_ENV=production"

[Install]
WantedBy=multi-user.target
